# src/backend/mysqlRollBack.cron
30 22 * * * /var/lib/mysql-files/backup/shell/mysql_dump_script.sh
